# 09-server-python/supervisord.conf
[program:flask_web_socket_app]
user=root                                                             ; To-run-as user
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log        ; Where to write log messages
stdout_logfile_maxbytes=10MB
redirect_stderr=true                                                  ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding
process_name=%(program_name)s_%(process_num)02d
command=sh /home/app/web/run.sh  

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log        ; Where to write log messages
stdout_logfile_maxbytes=10MB

[supervisord]
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
umask=022

[supervisorctl]
prompt=mysupervisor
